# board/conservify/glacier/overlay/etc/cron.hourly/cleanup
#!/bin/bash
# -*- mode: sh -*-

find /data/geophone/ -name "*.csv" -delete
