# board/conservify/glacier/overlay/etc/cron.hourly/morningstar
#!/bin/bash
# -*- mode: sh -*-

/usr/local/bin/morningstar --driver pl2303 --syslog morningstar > /dev/null
