# board/conservify/common/overlay/etc/cron.d/status
# Show system information every 5 minutes.

*/5 * * * * /usr/local/bin/status | /usr/bin/logger -t status
