# package/fake-hwclock/fake-hwclock.cron
*/15 * * * * /bin/fake-hwclock | /usr/bin/logger -t fake-hwclock
